# os/overlay/etc/systemd/user/ghost-daemon.service
[Unit]
Description=GhOSt system daemon (ghostd)

[Service]
ExecStart=/usr/local/bin/ghostd --listen 127.0.0.1:7700 --token-file %h/.config/ghost/token --static /usr/share/ghost/shell
Restart=always
RestartSec=1
# ghostd must never run as root; privileged actions go through NM/logind/polkit.
NoNewPrivileges=yes

[Install]
WantedBy=default.target
